# systemd/adblock-dashboard.service
[Unit]
Description=Ad-Block Web Dashboard
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/adblock
ExecStart=/opt/adblock/venv/bin/python3 dashboard_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
